# Build/Dockerfile/Debian/systemd/services/novnc.service
[Unit]
Description=Websockify noVNC Service
After=network.target

[Service]
User=root
ExecStart=/usr/bin/websockify -D --web /usr/share/novnc 6080 localhost:5901
Restart=always
RestartSec=5
#StartLimitInterval=5min
#StartLimitBurst=10

[Install]
WantedBy=multi-user.target
