# Build/Dockerfile/Debian/systemd/services/novnc.service
[Unit]
Description=NoVNC Server
After=network.target

[Service]
User=root
ExecStart=/usr/share/novnc 6080
#ExecStart=/usr/bin/novnc --listen 8080 --vnc localhost:5900
Restart=always

[Install]
WantedBy=multi-user.target
