# Build/Dockerfile/Debian/systemd/services/tightvnc.service
[Unit]
Description=TightVNC Server
After=network.target

[Service]
User=root
ExecStart=/usr/bin/tightvncserver :1
Restart=always
RestartSec=5
#StartLimitInterval=5min
#StartLimitBurst=10

[Install]
WantedBy=multi-user.target
