# Build/Dockerfile/Debian/systemd/services/websockify.service
[Unit]
Description=Websockify noVNC Service
After=network.target

[Service]
User=root
ExecStart=/usr/bin/websockify -D --web=/usr/share/novnc/ 6080 localhost:5900
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
