# Build/Dockerfile/Debian/systemd/myservice.service
[Unit]
Description=My Service
After=network.target

[Service]
User=myservice
ExecStart=/app/myservice.sh
Restart=always

[Install]
WantedBy=multi-user.target
